# gratia-probe/osg/gratia-probe-glexec.cron
0 * * * * root perl -e 'sleep rand 3600' && /usr/share/gratia/condor/condor_meter.cron.sh
