# gratia-probe/osg/gratia-probe-sge.cron
0,15,30,45 * * * * root perl -e 'sleep rand 900' && /usr/share/gratia/sge/sge_meter.cron.sh
