# gratia-probe/osg/gratia-probe-condor-events.cron
0,10,20,30,40,50 * * * * root perl -e 'sleep rand 600' && /usr/share/gratia/condor-events/watchCondorEvents
