# gratia-probe/osg/gratia-probe-dcache-storage.cron
0 * * * * root perl -e 'sleep rand 3600' && /usr/share/gratia/dCache-storage/dCache-storage_meter.cron.sh
