# gratia-probe/osg/gratia-probe-bdii-cese.cron
0,15,30,45 * * * * root perl -e 'sleep rand 900' && /usr/share/gratia/bdii-status/bdii_cese_record
