# gratia-probe/osg/gratia-probe-hadoop-storage.cron
0 * * * * root perl -e 'sleep rand 3600' && /usr/share/gratia/hadoop-storage/hadoop_storage_probe -c /etc/gratia/hadoop-storage/storage.cfg
