# gratia-probe/osg/gratia-probe-psacct.cron
0 1 * * * root perl -e 'sleep rand 9600' && /usr/share/gratia/psacct/psacct_probe.cron.sh
