# gratia-probe/osg/gratia-probe-gridftp-transfer.cron
0,30 * * * * root perl -e 'sleep rand 1800' && /usr/share/gratia/gridftp-transfer/gridftp-transfer_meter.cron.sh
